# ftp-daemon/soltrace-daemon.service
[Unit]
Description=SolTrace FTP Log Daemon
After=network-online.target proftpd.service
Wants=network-online.target
PartOf=proftpd.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=soltrace
Group=soltrace
WorkingDirectory=/opt/soltrace-daemon
ExecStart=/opt/soltrace-daemon/venv/bin/python3 soltrace_daemon.py

Restart=on-failure
RestartSec=10
RestartPreventExitStatus=1

TimeoutStartSec=30
TimeoutStopSec=15

# ── 보안 하드닝 (Rocky 8 / systemd 239+) ────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/soltrace /var/log /opt/soltrace-daemon
# /usr/service 상위 경로도 포함해야 ProtectSystem=strict 하에서 접근 가능
ReadOnlyPaths=/usr/service
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=soltrace-daemon

MemoryMax=64M
CPUQuota=10%
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
